# go/contrib/kyuko.service
[Unit]
Description=Cancelled class notificator for Doshisha University

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kyuko
